# mktk-provisioner.service
[Unit]
Description=Mikrotik RouterBoard Provisioner
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/mktk-provisioner
Environment=FLASK_APP=app.py
Environment=FLASK_ENV=production
ExecStart=/home/pi/mktk-provisioner/venv/bin/python /home/pi/mktk-provisioner/app.py
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/mktk-provisioner

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mktk-provisioner

[Install]
WantedBy=multi-user.target
